# Docker/prod/app/supervisor/messanger.conf
[program:intranet-consumer]
command=/var/www/html/app/bin/console messenger:consume
directory=/var/www/html/app
user=www-data
autostart=true
autorestart=true
exitcodes=0,10
numprocs=1
startretries=5
